# contrib/systemd/astral-net-autotune@.service
[Unit]
Description=ASTRAL Network Autotune (%i)
After=network.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/home/hex/astra
EnvironmentFile=/etc/astral/%i.env
# Optional overrides (can be placed into /etc/astral/%i.env):
# AUTOTUNE_USERNAME=admin
# AUTOTUNE_PASSWORD=admin
# AUTOTUNE_DURATION_SEC=240
# AUTOTUNE_CANDIDATES=bad,max,superbad
# AUTOTUNE_JITTER_VARIANTS_MS=2000,6000,12000
TimeoutStartSec=3600
ExecStart=/bin/sh -lc 'exec python3 -u /home/hex/astra/tools/net_autotune.py --api "http://127.0.0.1:${PORT}" --username "${AUTOTUNE_USERNAME:-admin}" --password "${AUTOTUNE_PASSWORD:-admin}" --duration-sec "${AUTOTUNE_DURATION_SEC:-240}" --candidates "${AUTOTUNE_CANDIDATES:-bad,max,superbad}" --jitter-variants-ms "${AUTOTUNE_JITTER_VARIANTS_MS:-2000,6000,12000}"'

[Install]
WantedBy=multi-user.target
