# contrib/systemd/stream-watchdog.service
[Unit]
Description=Stream watchdog (CPU/RAM) for Stream process
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/stream-watchdog.sh
