# contrib/systemd/stream-net-autotune@.timer
[Unit]
Description=Stream Network Autotune Timer (%i)

[Timer]
# Автотюнинг занимает 3-5 минут на стрим и делает изменения только для входов с `#net_tune=1`.
# Запускаем не слишком часто и с рандомной задержкой, чтобы не "долбить" источники по расписанию.
OnCalendar=*-*-* 00,06,12,18:25:00
RandomizedDelaySec=1800
Persistent=true
Unit=stream-net-autotune@%i.service

[Install]
WantedBy=timers.target
